# deploy/sunshine-backup.service
[Unit]
Description=Backup SunShine SQLite database

[Service]
Type=oneshot
ExecStart=/opt/sunshine/my-backend/deploy/backup-sqlite.sh

